# rootfs/etc/services.d/blackout/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting Blackout Time Parser"
cd /opt/blackout-time-parser
exec node dist/index.js
